# site-cookbooks/brewbit/templates/default/etc/cron.daily/backup.erb
#!/bin/sh
su - deploy -c "cd /var/www/<%= node[:brewbit][:hostname] %>/current/backup/ && bundle install && RAILS_ENV=<%= node['environment'] %> backup perform --root-path . --trigger db"
