# site-cookbooks/brewbit/files/default/etc/cron.daily/backup
#!/bin/sh
su - deploy -c "cd /var/www/brewbit.com/current/backup && bundle install && RAILS_ENV=production backup perform --root-path . --trigger db
